# deploy/services/fastapi.service
[Unit]
Description=FastAPI server

[Service]
User=pi
WorkingDirectory=/home/pi/projects/reddit
Environment=poetry=/home/pi/.poetry/bin/poetry
ExecStart=export $(grep -v '^#' .env | xargs) && poetry shell && newrelic-admin run-program poetry run fastapi
Restart=always

[Install]
WantedBy=multi-user.target
